# src/linux/debian/prerm
#!/bin/sh

if [ "$1" = "remove" ]; then
    service goldilocks stop
    update-rc.d -f goldilocks remove
fi
